# deploy/openfang.service
[Unit]
Description=ClawForge Agent OS Daemon
Documentation=https://github.com/clawforge-ai/clawforge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=clawforge
Group=clawforge
ExecStart=/usr/local/bin/clawforge start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
EnvironmentFile=-/etc/clawforge/env
WorkingDirectory=/var/lib/clawforge

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/clawforge
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
